# containers/.images/php5-cron/conf/cron.d/magicpro
* *     * * *     www-data   /usr/bin/php /var/www/cli.php >> /var/www/vardata/log/crontab/cron.log 2>&1
* *     * * *     cuca   /usr/bin/php /var/www/cli.php crontab:scheduler@update >> /var/www/vardata/log/crontab/update.log 2>&1
