# containers/images/apache2-php5/conf/cron.d/magicpro
* *     * * *     www-data   /usr/bin/php /var/www/index.php >> /var/www/vardata/log/crontab/cron.log 2>&1
